# shell/profile.d/node.profile
NODE_PATH=`brew --prefix`/lib/node

export PATH=${PATH}:${NODE_PATH}
