# shell/profile.d/rbx.profile
echo "Rubinius...."
export RBXOPT="-Xrbc.db"
