# bash/profile.d/env.profile
##
# General User environment
export EDITOR=vim
export PAGER=most

export PATH=${HOME}/bin:${PATH}
