# bash/profile.d/env.profile
##
# General User environment
export EDITOR=vim
export PAGER=most
export VIM_APP_DIR=/usr/local/Cellar/macvim/latest

export PATH=${HOME}/bin:${PATH}