# bash/profile.d/env.profile
##
# General User environment
export EDITOR=vim
export PAGER=most

##
# Gems stuff
export GEM_PRIVATE_KEY='~/config/pki/keys/gem-private_key.pem'
export GEM_CERTIFICATE_CHAIN='~/config/pki/certs/gem-public_cert.pem'


##
# Path
export PATH=/usr/local/bin:/usr/local/sbin:/usr/local/mysql/bin:~/bin:${PATH}
