# bash/profile.d/env.profile
##
# General User environment
export EDITOR=mvim
export PAGER=most
export VIM_APP_DIR=/usr/local/Cellar/macvim/v7.3-53

export PATH=${HOME}/bin:${PATH}
