# bash/profile.d/env.profile
##
# General User environment
export EDITOR=mvim
export PAGER=most

export PATH=${HOME}/bin:${PATH}