# bash/profile.d/gem.profile
##
# Gems stuff
export GEM_PRIVATE_KEY='${HOME}/config/pki/keys/gem-private_key.pem'
export GEM_CERTIFICATE_CHAIN='${HOME}/config/pki/certs/gem-public_cert.pem'