# bash/profile.d/node.profile
NPM_HOME=/usr/local/share/npm

export PATH=${PATH}:${NPM_HOME}/bin